# daemon/claude-usage-daemon.service
[Unit]
Description=Clawdmeter Usage Daemon (HTTP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=__VENV_PYTHON__ __DAEMON_PATH__
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
